# services/pyha.socket
[Unit]
Description=pyha socket

[Socket]
ListenStream=/run/pyha/socket
ListenStream=
ListenStream=

[Install]
WantedBy=sockets.target